# cec_mqtt_bridge/rootfs/etc/services.d/cec_mqtt_bridge/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting HDMI-CEC MQTT Bridge"

export PYTHONPATH="/app/src"
exec python3 -u -m cec_mqtt_bridge.bridge
